# scripts/systemd/sujbot-backup.timer.template
[Unit]
Description=SUJBOT Daily Backup Timer
Documentation=https://github.com/ADS-teamA/SUJBOT

[Timer]
# Run daily at 3:00 AM
OnCalendar=*-*-* 03:00:00
# Randomize start time by up to 15 minutes to avoid thundering herd
RandomizedDelaySec=900
# Run immediately if last scheduled run was missed
Persistent=true

[Install]
WantedBy=timers.target
