# scripts/systemd/sujbot-backup.service
[Unit]
Description=SUJBOT Database Backup
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=michal
WorkingDirectory=/home/michal/SUJBOT
ExecStart=/home/michal/SUJBOT/scripts/backup_cron.sh
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/michal/SUJBOT/backups /home/michal/SUJBOT/logs

[Install]
WantedBy=multi-user.target
